# .data/topgrade.service.tmpl
# Systemd User Service for Topgrade
# Purpose: Automated system maintenance using topgrade
# Template Variables: Uses chezmoi template system for user-specific configuration

[Unit]
Description=Topgrade System Maintenance
Documentation=https://github.com/topgrade-rs/topgrade
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Run topgrade with configuration file
ExecStart=/usr/bin/topgrade --config {{ .chezmoi.homeDir }}/.config/topgrade.toml
# Set working directory to user home
WorkingDirectory={{ .chezmoi.homeDir }}
# User services automatically run as the correct user
# Environment variables
Environment=HOME={{ .chezmoi.homeDir }}
Environment=USER={{ .chezmoi.username }}
Environment=XDG_CONFIG_HOME={{ .chezmoi.homeDir }}/.config
Environment=XDG_DATA_HOME={{ .chezmoi.homeDir }}/.local/share
Environment=XDG_CACHE_HOME={{ .chezmoi.homeDir }}/.cache
# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=topgrade
# Security settings (relaxed for system maintenance)
NoNewPrivileges=false
PrivateTmp=false
ProtectSystem=false
ProtectHome=false
# Timeout settings
TimeoutStartSec=1800
TimeoutStopSec=300

[Install]
WantedBy=default.target
