# .data/topgrade.service.tmpl
# Systemd System Service for Topgrade
# Purpose: Automated system maintenance using dedicated maintenance user
# Template Variables: Uses chezmoi template system for configuration
# User: topgrade-maintenance (dedicated system user for security isolation)

[Unit]
Description=Topgrade System Maintenance
Documentation=https://github.com/topgrade-rs/topgrade
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
# Run as dedicated maintenance user for security isolation
User=topgrade-maintenance
Group=topgrade-maintenance
# Run topgrade with system configuration file
ExecStart=/usr/bin/topgrade --config /etc/topgrade/topgrade.toml
# Set working directory to maintenance user home
WorkingDirectory=/var/lib/topgrade-maintenance
# Environment variables for maintenance user
Environment=HOME=/var/lib/topgrade-maintenance
Environment=USER=topgrade-maintenance
Environment=XDG_CONFIG_HOME=/etc/topgrade
Environment=XDG_DATA_HOME=/var/lib/topgrade-maintenance/.local/share
Environment=XDG_CACHE_HOME=/var/lib/topgrade-maintenance/.cache
# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=topgrade
# Security settings (enhanced for maintenance user)
NoNewPrivileges=false
# Protect user home directories from maintenance user access
ProtectHome=yes
# Use private /tmp for isolation
PrivateTmp=yes
# Allow system modifications (needed for package updates)
ProtectSystem=false
# Timeout settings
TimeoutStartSec=1800
TimeoutStopSec=300

[Install]
WantedBy=multi-user.target
