# private_dot_config/systemd/user/voxtype.service.d/gpu.conf.tmpl
{{ if eq .chassisType "laptop" }}
[Service]
Environment="VOXTYPE_VULKAN_DEVICE=nvidia"
{{ end }}
