# private_dot_config/systemd/user/home-backup.timer
[Unit]
Description=Daily home backup at 2am

[Timer]
OnCalendar=*-*-* 02:00:00
Persistent=true
RandomizedDelaySec=30min

[Install]
WantedBy=timers.target
