# private_dot_config/systemd/user/home-backup.service
[Unit]
Description=Incremental home backup (restic)
After=network.target

[Service]
Type=oneshot
# Only run when the backup repository is initialized
ExecCondition=/bin/sh -c 'test -d "%h/.local/share/backups/home-backup"'
ExecStart=%h/.local/lib/scripts/system/home-backup backup --non-interactive
ExecStart=%h/.local/lib/scripts/system/home-backup prune
Restart=no
TimeoutSec=3600
Environment=PATH=%h/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/bin
Environment=HOME=%h
Environment=RESTIC_REPOSITORY=%h/.local/share/backups/home-backup
Environment=RESTIC_PASSWORD_FILE=%h/.local/share/backups/.restic-password

[Install]
WantedBy=default.target
