# private_dot_config/systemd/user/hyprdynamicmonitors-prepare.service
[Unit]
Description=HyprDynamicMonitors - Prepare monitor config before Hyprland
Before=hyprland.service

[Service]
Type=oneshot
ExecStart=/usr/bin/hyprdynamicmonitors prepare
RemainAfterExit=yes

[Install]
WantedBy=hyprland.service
