# private_dot_config/systemd/user/llama-server.service.tmpl
[Unit]
Description=llama.cpp inference server
After=network-online.target
ConditionPathExists=%h/.config/llama-server/env

[Service]
Type=simple
EnvironmentFile=%h/.config/llama-server/env
ExecStart=/usr/bin/llama-server \
  --model %h/.local/share/models/${MODEL_FILE} \
  --host {{ .ai.server.host }} \
  --port {{ .ai.server.port }} \
  --n-gpu-layers ${N_GPU_LAYERS:-99}
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
