# private_dot_config/systemd/user/system-health-check.timer
[Unit]
Description=System health check timer (every 15 minutes)

[Timer]
OnBootSec=10min
OnUnitActiveSec=15min
Persistent=false
RandomizedDelaySec=1min

[Install]
WantedBy=timers.target
