# private_dot_config/systemd/user/wallpaper-cycle.service
[Unit]
Description=Automatic wallpaper rotation service
Documentation=man:swww(1)
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=%h/.config/scripts/media/random-wallpaper.sh fade 3
# Use fade transition with 3 second duration for smooth cycling

[Install]
WantedBy=default.target
