# private_dot_config/systemd/user/wallpaper-cycle.timer
[Unit]
Description=Automatic wallpaper rotation timer (every 30 minutes)
Documentation=man:systemd.timer(5)

[Timer]
# Trigger 5 minutes after boot (initial run), then every 30 minutes after service completes
OnBootSec=5min
OnUnitActiveSec=30min
# Don't run immediately if a scheduled run was missed
Persistent=false
# Small random delay to prevent exact timing conflicts
RandomizedDelaySec=2min

[Install]
WantedBy=timers.target
